# scripts/sailframes-sync.timer
[Unit]
Description=Run SailFrames sync periodically
Documentation=https://github.com/sailframes/sailframes

[Timer]
# Run 2 minutes after boot
OnBootSec=2min

# Run every 30 minutes
OnUnitActiveSec=30min

# Add random delay up to 5 minutes to avoid thundering herd
RandomizedDelaySec=5min

# Persist timer across reboots
Persistent=true

[Install]
WantedBy=timers.target
