# scripts/sailframes-sync.service
[Unit]
Description=SailFrames Data Sync to AWS S3
Documentation=https://github.com/sailframes/sailframes
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sailframes-sync.sh
User=root
Group=root

# Environment
Environment=SAILFRAMES_CONFIG=/etc/sailframes/sailframes.yaml
Environment=SAILFRAMES_DATA_DIR=/mnt/sailframes-data
Environment=SAILFRAMES_S3_BUCKET=sailframes-fleet-data-prod

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sailframes-sync

# Retry on failure
Restart=on-failure
RestartSec=60

# Resource limits
Nice=10
IOSchedulingClass=idle

[Install]
WantedBy=multi-user.target
